# Raspberry Pi 5/services/cpu_logger.service
[Unit]
Description=CPU Activity Logging
After=network.target

[Service]
User=nas
WorkingDirectory=/home/nas/services/
ExecStart=/usr/bin/python /home/nas/services/cpu_logger.py

[Install]
WantedBy=multi-user.target
